# packaging/postinstall.sh
#!/bin/sh
echo "GopherClaw installed. To enable auto-start:"
echo "  systemctl --user daemon-reload"
echo "  systemctl --user enable --now gopherclaw"
